# deploy/nextapp.service
[Unit]
Description=Next.js SSR App
After=network.target

[Service]
Type=simple
WorkingDirectory=/srv/nextapp
ExecStart=/usr/bin/npm run start -- -p 3000
Restart=always
Environment=NODE_ENV=production
# Opcional: límites de memoria y tiempo de parada
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
