# deploy/sommelier_backend.service
[Unit]
Description=Sommelier Backend Service
After=network.target

[Service]
# Usuario y grupo no-root
User=sommelier
Group=sommelier

# Directorio de trabajo
WorkingDirectory=/srv/sommelier_backend

# Variables de entorno (se reemplazan en el workflow con sed)
Environment=OPENAI_API_KEY=changeme
Environment=HEYGEN_API_KEY=changeme

# Activar entorno virtual
ExecStart=/srv/sommelier_backend/venv/bin/gunicorn \
          --workers 3 \
          --bind 0.0.0.0:5000 \
          app:app

# Reinicio automático en caso de fallo
Restart=always
RestartSec=5

# Logs centralizados en journalctl
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
