# deploy/sommelier_backend.service
[Unit]
Description=Sommelier Backend (Flask API)
After=network.target

[Service]
Type=simple
WorkingDirectory=/srv/sommelier_backend
ExecStart=/usr/bin/python3 app.py
Restart=always
RestartSec=5
Environment=FLASK_ENV=production
Environment=PORT=5000
Environment=CHROMADB_DISABLE_TELEMETRY=1
Environment=CHROMADB_PATH=./chroma_db
Environment=OPENAI_API_KEY=${OPENAI_API_KEY}
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
